# meta-ismail/recipes-applications/hello_cmake/files/hello-startup.service
[Unit]
Description=Hello CMake Service

[Service]
ExecStart=/usr/bin/hello-cmake

[Install]
WantedBy=multi-user.target