# meta-ismail/recipes-applications/vehicle/files/vehicle-startup.service
[Unit]
Description=vehicle CMake Service

[Service]
ExecStart=/usr/bin/vehicle_diagnostics 

[Install]
WantedBy=multi-user.target